# etc/init.d/xow
#!/sbin/openrc-run

command="/usr/bin/xow"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"
