# etc/init.d/zerotier-one
#!/sbin/openrc-run

depend() {
    provide net
    use network
    before firewall netmount 
}

command="/usr/sbin/zerotier-one"
#command_args="-U"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"
