# etc/init.d/zerotier-one
#!/sbin/openrc-run

depend() {
    need net
    use network
#    before firewall netmount
#    after connmand
}

command="/usr/sbin/zerotier-one"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"
